# files/systemd/apps-backup.timer
[Unit]
Description=RSync of the Apps to Tank for backup daily at 11 PM

[Timer]
OnCalendar=*-*-* 23:00:00

[Install]
WantedBy=timers.target
