# files/systemd/duck.timer
[Unit]
Description=DuckDNS IP update every 5 minutes.

[Timer]
OnActiveSec=5m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
