# files/systemd/zfs-snapshots.timer
[Unit]
Description=ZFS snapshots daily at midnight

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
